# BCINFO3/root/.crontab
# email output to me, not locally
# TODO: this wont work until I set up SOME sort of SMTP server on bcinfo3
# removing this, hoping to have email go to /var/spool/mail/root or
# something which I can then download remotely
# MAILTO="bcinfo3@barrycarter.info"
# reload self every hour (meaning I can just update this file instead
# of having to run "crontab -e"; however, I DO have to run crontab
# ~/.crontab once)
0 * * * * crontab /root/.crontab
# update http://test.barrycarter.info/sunstuff.html
* * * * * /usr/local/bin/bc-sun-always-shines.pl
# update http://cheapbay.barrycarter.info/
*/5 * * * * /usr/local/bin/bc-cheap-bay.pl
# send email to comicssherpa.com authors who have requested it
*/10 * * * * r2e run
# keeps /tmp clean (I tend to create a lot of /tmp files)
42 1 * * * find /tmp -name 'cache*' -mmin +180 | xargs --no-run-if-empty rm
# update database of files for locate
0 2 * * * updatedb
# remove cache files
28 * * * * find /var/tmp/cache -type f -mmin +60 | xargs rm
# and radar
04 3 * * * find /var/tmp/radar -mtime +1 | xargs rm
# this is really ugly (TODO: clean this up)
54 4 * * * find /sites/data/ -mtime +1 -name '[A-Z][A-Z][A-Z]*_N??.gif' | xargs rm
59 5 * * * find /sites/data/ -mtime +1 -name 'Conus*.gif' | xargs rm
# TODO: below is really ugly, but avoids find command
09 * * * * rm /var/tmp/querys/DONE/*
# TODO: should really use tmpwatch here or access time instead of mtime...
09 6 * * * find /tmp/ -mtime +7 | xargs rm
# MOS guidance
52 * * * * /usr/local/bin/bc-get-guidance2.pl
# check that all daemons are running properly
*/5 * * * * /usr/local/bin/bc-daemon-checker2.pl --sleep=50
# keep geonames changes up to date (hourly here is overkill but ok)
3 * * * * /usr/local/bin/bc-daily-update.pl --debug
# dump gocomics db nightly (only one I actually create remotely)
7 1 * * * cd /sites/DB/; /usr/local/bin/sqlite3 gocomics.db .dump > gocomics-dump.txt; mv gocomics-dump.txt.bz2 gocomics-dump.txt.bz2.old; bzip2 gocomics-dump.txt
